# docker/images/rpi-spark-slave/spark-slave.conf
[program:spark-slave]
command=/bin/bash -c "/opt/spark/sbin/start-slave.sh spark://pi1-docker:7077"
redirect_stderr=true

