# docker/images/rpi-spark-slave/spark-slave.conf
[program:spark-slave]
command=/bin/bash -c "export SPARK_PUBLIC_DNS=$HOSTNAME && /opt/spark/sbin/start-slave.sh $SPARK_MASTER"
redirect_stderr=true

