# docker/mpi/base/sshd.conf
[program:sshd]
command=/bin/bash -c "/etc/init.d/ssh start"
redirect_stderr=true

