# docker/images/rpi-spark-master/spark-master.conf
[program:spark-master]
command=/bin/bash -c "/opt/spark/sbin/start-master.sh"
redirect_stderr=true

