# docker/images/plattforms/x86/rpi-spark-slave/spark-slave.conf
[program:spark-slave]
command=/bin/bash -c "/opt/spark/sbin/start-slave.sh --host $HOSTNAME $SPARK_MASTER"
redirect_stderr=true

